# itmo-553/week-09/service-files/graphiteb/graphite-api.service
[Unit]
Description=Graphite-API service
Requires=graphite-api.socket
RequiresMountsFor=/var/log/graphite-api /var/lib/graphite/whisper /var/lib/graphite-api
ConditionFileIsExecutable=/usr/bin/gunicorn
ConditionPathExists=/etc/graphite-api.yaml
ConditionPathExists=/home/vagrant/.local/lib/python3.6/site-packages/graphite_api/app.py

[Service]
ExecStart=/usr/bin/gunicorn --name=graphite-api graphite_api.app:app
Restart=on-failure
User=_graphite
Group=_graphite
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s TERM $MAINPID
PrivateTmp=true

[Install]
WantedBy=multi-user.target